# data/udev/99-fanatec-wheel-perms.rules
# Match kernel name of device, rather than ATTRS{idProduct} and ATTRS{idVendor}
# so we can access the range file and leds directory. Set rw access to these 
# files for everyone.

# FANATEC CSL Elite Wheel Base
SUBSYSTEMS=="hid", KERNELS=="0003:0EB7:0E03.????", DRIVERS=="ftec_csl_elite", RUN+="/bin/sh -c 'cd %S%p/../../../; chmod 666 range'"

# FANATEC CSL Elite Wheel Base PS44
SUBSYSTEMS=="hid", KERNELS=="0003:0EB7:0005.????", DRIVERS=="ftec_csl_elite", RUN+="/bin/sh -c 'cd %S%p/../../../; chmod 666 range'"

# FANATEC ClubSport Wheel Base V2
SUBSYSTEMS=="hid", KERNELS=="0003:0EB7:0001.????", DRIVERS=="ftec_csl_elite", RUN+="/bin/sh -c 'cd %S%p/../../../; chmod 666 range'"

# FANATEC ClubSport Wheel Base V2.5
SUBSYSTEMS=="hid", KERNELS=="0003:0EB7:0004.????", DRIVERS=="ftec_csl_elite", RUN+="/bin/sh -c 'cd %S%p/../../../; chmod 666 range'"

# FANATEC Podium Wheel Base DD1
SUBSYSTEMS=="hid", KERNELS=="0003:0EB7:0006.????", DRIVERS=="ftec_csl_elite", RUN+="/bin/sh -c 'cd %S%p/../../../; chmod 666 range'"
